# linkyd.service
[Service]
Type=simple
ExecStart=/opt/linkyd

[Install]
WantedBy=multi-user.target
