# linkyd.service
[Service]
Type=simple
ExecStart=/home/feo/linky-start.sh

[Install]
WantedBy=multi-user.target
